# package/DEBIAN/postinst
#!/bin/sh
set -e

systemctl daemon-reload

systemctl enable ndm-test-server.service

systemctl start ndm-test-server.service